# apps/facebook-token-cloak/launchd/com.affiliate.facebook-token-cloak.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.affiliate.facebook-token-cloak</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>HOME=/Users/yok-macmini</string>
    <string>NODE_ENV=production</string>
    <string>PORT=8820</string>
    <string>node</string>
    <string>/Users/yok-macmini/Developer/shopping-affiliate/apps/facebook-token-cloak/bin/start.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/yok-macmini/Developer/shopping-affiliate/apps/facebook-token-cloak</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/yok-macmini</string>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>PATH</key>
    <string>/Users/yok-macmini/.local/opt/node-v24.14.1-darwin-arm64/bin:/Users/yok-macmini/.local/opt/node/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <!-- Accounts Bridge v2 default account. Must match the configured Facebook Lite/Power Editor role owner; do not point this at a legacy profile. -->
    <key>FACEBOOK_TOKEN_CLOAK_POST_ACCOUNT</key>
    <string>100090320823561</string>
    <!-- Non-secret Ads Manager TEMPLATE_SALES adset id. Keep runtime explicit so the bridge never
         falls back to the retired pre-SALES template adset. -->
    <key>FACEBOOK_TOKEN_CLOAK_TEMPLATE_ADSET</key>
    <string>120248134990230263</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/yok-macmini/Library/Logs/facebook-token-cloak.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/yok-macmini/Library/Logs/facebook-token-cloak.err.log</string>
</dict>
</plist>
