# apps/affiliate-credential-vault/launchd/com.affiliate.credential-vault.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict><key>Label</key><string>com.affiliate.credential-vault</string><key>ProgramArguments</key><array><string>/bin/zsh</string><string>-lc</string><string>cd /Users/yok-macmini/Developer/shopping-affiliate/apps/affiliate-credential-vault &amp;&amp; HOME=/Users/yok-macmini npm start</string></array><key>RunAtLoad</key><true/><key>KeepAlive</key><true/><key>StandardOutPath</key><string>/Users/yok-macmini/Library/Logs/affiliate-credential-vault.log</string><key>StandardErrorPath</key><string>/Users/yok-macmini/Library/Logs/affiliate-credential-vault.err.log</string></dict></plist>
